# distro/x86_64-live/live-build/config/includes.chroot/etc/systemd/system/geos-telemetry.service
[Unit]
Description=GeOS Telemetry Collector
After=geos-energy-controller.service geos-workload-manager.service
Requires=geos-energy-controller.service

[Service]
Type=simple
WorkingDirectory=/opt/geos
Environment=PYTHONPATH=/opt/geos
ExecStart=/opt/geos/.venv/bin/python -m telemetry.collector
Restart=always
RestartSec=2

[Install]
WantedBy=geos-demo.target
