# distro/x86_64-live/live-build/config/includes.chroot/etc/xdg/autostart/geos-gui.desktop
[Desktop Entry]
Type=Application
Name=GeOS
Comment=Launch the GeOS demo shell
Exec=/usr/local/bin/geos-launch-gui
X-GNOME-Autostart-enabled=true
Terminal=false
